# tools/shadow/launchd/com.zeta.shadow-observer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.zeta.shadow-observer.plist — launchd LaunchAgent template for the Zeta
  shadow observer (B-0402 Slice 6 — durable tick-source via macOS launchd).

  This file is a TEMPLATE. Mustache-style placeholders in the body get
  replaced with the local bun path and the local Zeta repo root.

  DO NOT manually substitute. Run the installer:

      bun tools/shadow/launchd/install-launchagent.ts

  The installer substitutes the placeholders safely (no shell-metacharacter
  risk with paths containing &, |, \), runs plutil -lint on the result,
  writes to ~/Library/LaunchAgents/, and optionally bootstraps via
  launchctl. See tools/shadow/launchd/README.md for the full install
  procedure (incl. accessibility permission) and the uninstall recipe.

  Composes-with:
    - tools/shadow/shadow-observer.ts (the process this launches)
    - docs/backlog/P0/B-0402-zeta-shadow-mode-first-class-cli-autocomplete.md
    - .claude/rules/shadow-star-shorthand-autocomplete-marker.md (the marker
      this observer-loop produces when accepted suggestions get shipped)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zeta.shadow-observer</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{BUN_PATH}}</string>
        <string>{{REPO_ROOT}}/tools/shadow/shadow-observer.ts</string>
        <string>--loop</string>
        <string>5000</string>
        <string>--loop-interval</string>
        <string>2000</string>
        <string>--delay</string>
        <string>3000</string>
        <string>--dry-run</string>
        <string>--log-file</string>
        <string>{{REPO_ROOT}}/tools/shadow/shadow-observer.log</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Crash-only restart: KeepAlive dict keys are OR-combined.
             SuccessfulExit=false would also relaunch on any non-zero exit
             (configuration errors, missing accessibility permission), which
             would mask config bugs as crash-loops. Crashed=true on its own
             gives the documented "relaunch on signal-termination only"
             semantic. Per Codex review on PR #3342. -->
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>{{REPO_ROOT}}/tools/shadow/shadow-observer.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{REPO_ROOT}}/tools/shadow/shadow-observer.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
